# server/fog-vizu.service
[Unit]
Description=Fog Gate Visualizer Backend Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/fog-vizu/server
ExecStart=/var/www/fog-vizu/server/venv/bin/uvicorn fogvizu.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
EnvironmentFile=/var/www/fog-vizu/server/.env

[Install]
WantedBy=multi-user.target
